# ops/orien-monitor.cron
*/5 * * * * root /srv/apps/orien_comercial/app/ops/monitor-health.sh >/dev/null 2>&1
